# root/etc/s6-overlay/s6-rc.d/init-targetcli/run
#!/usr/bin/with-contenv bash

targetctl restore
